# 01_pyTREMOR_lights/pytremor_lights.service
[Unit]
Description=pyTREMOR Seismic LED Driver (live volcanic tremor)
After=network-online.target pigpiod.service
Wants=network-online.target
Requires=pigpiod.service

[Service]
Type=simple
User=sjc1
WorkingDirectory=/home/sjc1
ExecStart=/home/sjc1/venv_tremor/bin/python3 /home/sjc1/pyTREMOR_lights01.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/pytremor_lights.log
StandardError=append:/var/log/pytremor_lights.log

[Install]
WantedBy=multi-user.target
